# sysfiles/udev-rules/user-btkeyboard.rules
# Set bluetooth keyboard power up
# written 2019-07-18
ACTION=="add", KERNEL=="hci0", RUN+="/bin/bluetooth --agent menu power on"
